# docker-murmur/services.d/murmur/run
#!/bin/sh
exec /opt/murmur/murmur.x86 -fg -v -ini /etc/murmur.ini
