# murmur/services.d/gtmurmur/run
#!/bin/sh
exec /usr/local/bin/gtmurmur-static -conf /etc/murmur.ini
